# blog.ini
[program:blog]
command=/home/shijl0925/May-Blog/venv/bin/gunicorn wsgi:app -c gconfig.py
directory=/home/shijl0925/May-Blog/
stdout_logfile=/home/shijl0925/log/blog.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/home/shijl0925/log/blog.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true