# supervisord.d/blog.ini
[program:blog]
command=/usr/bin/python3 /usr/local/bin/gunicorn wsgi:app -c gconfig.py --bind 0.0.0.0:8080
directory=/Blog/
stdout_logfile=/home/logs/blog.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/home/logs/blog.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true